# systemd/alphaos-vault-sync-pull.service
[Unit]
Description=AlphaOS Vault pull on boot (remote -> local)
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
Environment="VAULT_REMOTE=fabian:AlphaOS-Vault"
Environment="VAULT_REMOTE_BACKUP=fabian:AlphaOS-Vault-backups"
EnvironmentFile=-/etc/alphaos-hub/env
EnvironmentFile=-/etc/alphaos-hub/alphaos-vault-sync.env
ExecStart=/usr/bin/bash -lc 'cd "$AOS_HUB_DIR/index-node" && exec ./scripts/vault-sync.sh pull'
User=alpha

[Install]
WantedBy=multi-user.target
