# systemd/aos-vault-push-eldanioo.service
[Unit]
Description=αOS Vault push to eldanioo (GDrive)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/bin/bash -lc 'rclone copy --update --skip-links --create-empty-src-dirs --backup-dir "eldanioo:MeineAblage/vault-backups/push-$(date +%%Y%%m%%d-%%H%%M%%S)" "%h/Dokumente/vault" "eldanioo:MeineAblage/vault"'
