# systemd/aos-vault-push-eldanioo.timer
[Unit]
Description=Daily αOS Vault push to eldanioo

[Timer]
OnCalendar=*-*-* 03:30:00
Persistent=true
Unit=aos-vault-push-eldanioo.service

[Install]
WantedBy=timers.target
