# systemd/aos-sync@.timer
[Unit]
Description=αOS Sync Timer (%i)

[Timer]
OnCalendar=*-*-* 23:00:00
Persistent=true

[Install]
WantedBy=timers.target
