# systemd/aos-hub-push.timer
[Unit]
Description=Daily 23:00 aos-hub push
Requires=aos-hub-push.service

[Timer]
OnCalendar=*-*-* 23:00:00
Persistent=true

[Install]
WantedBy=timers.target
