# systemd/aos-index.service
[Unit]
Description=AOS Hub - Index Node (HTTP server for menu.yaml)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=alpha
Group=alpha
EnvironmentFile=-/etc/alphaos-hub/env
EnvironmentFile=-/etc/alphaos-hub/aos.env
EnvironmentFile=-/etc/alphaos-hub/alphaos-index.env
ExecStart=/usr/bin/bash -lc 'cd "$AOS_HUB_DIR/index-node" && exec /usr/bin/node server.js'
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
