# systemd/core4-seed-week.timer
[Unit]
Description=Run Core4 seed-week every Monday morning

[Timer]
OnCalendar=Mon *-*-* 06:00:00
Persistent=true

[Install]
WantedBy=timers.target
