# systemd/core4-month-close.service
[Unit]
Description=Core4 month close (write monthly CSV for previous month)
ConditionPathExists=%h/aos-hub/core4/core4-month-close

[Service]
Type=oneshot
ExecStart=%h/aos-hub/core4/core4-month-close
