# systemd/user/aos-pwa-dev.service
[Unit]
Description=αOS Standalone PWA Runtime (Dev/Normal via pwa.env)
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/aos-hub/index-node
EnvironmentFile=-%h/.env/aos.env
EnvironmentFile=-%h/.env/pwa.env
Environment=NODE_ENV=development
Environment=PWA_HOST=0.0.0.0
Environment=PWA_PORT=8780
ExecStart=/usr/bin/bash -lc 'set -euo pipefail; mode="${AOS_PWA_RUN_MODE:-dev}"; case "$mode" in dev|development) exec /usr/bin/npm run pwa:dev ;; normal|prod|production) exec /usr/bin/npm run pwa ;; *) echo "invalid AOS_PWA_RUN_MODE=$mode (use dev|normal)" >&2; exit 2 ;; esac'
Restart=on-failure
RestartSec=1
CapabilityBoundingSet=
AmbientCapabilities=

StandardOutput=journal
StandardError=journal

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only

ReadOnlyPaths=%h/aos-hub/index-node
ReadWritePaths=%h/.task
ReadWritePaths=%h/.local/share/alphaos
ReadWritePaths=%h/.aos
ReadWritePaths=%h/vault

ProtectControlGroups=true
ProtectKernelTunables=true
ProtectKernelModules=true
LockPersonality=true
RestrictSUIDSGID=true

[Install]
WantedBy=default.target
