# systemd/user/aos-bridge.service
[Unit]
Description=AOS Bridge (aiohttp)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=-%h/.env/aos.env
WorkingDirectory=%h/aos-hub/bridge
ExecStart=/usr/bin/python %h/aos-hub/bridge/app.py --host 0.0.0.0 --port 8080
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
